# service/product-runtime-v34-3.service.template
[Unit]
Description=Product Runtime v34.3 Preproduction Worker fn01
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{RUNTIME_USER}}
WorkingDirectory={{RUNTIME_ROOT}}/product
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
ExecStart={{RUNTIME_PYTHON}} {{RUNTIME_ROOT}}/product/src/sync_xserver_main.py
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
